# build/supervised-asterisk.conf
[supervisord]
nodaemon=true

[group:asterisk]
programs=asterisk

[program:asterisk]
command= /usr/sbin/asterisk -g -f -p -U asterisk
autostart=true
autorestart=true
user=asterisk
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
priority=1

